# crontab
*/2 * * * * python /srv/unwatermark/erase.py
