# server-configuration/crontab
*/2 * * * * ps aux|grep -v grep|grep -E '^fadelee.* scp .*$' || python /srv/unwatermark/erase.py
